# config/supervisor.conf
[program:mcindoe]
directory=/home/ben/mcindoe
command=python manage.py run_gunicorn -c config/gunicorn.conf.py
user=root
autostart=true
autorestart=true
